# assets/supervisor/puma.conf
[program:puma]
priority=10
directory=/usr/local/share/tracks
command=bundle exec puma config.ru -e production -b unix:///tmp/tracks.socket
user=root
autostart=true
autorestart=true
stopsignal=QUIT
